# playbooks/roles/c2/files/papertrail.service
[Unit]
After=systemd-journald.service c2.service
Requires=systemd-journald.service c2.service

[Service]
ExecStart=/bin/sh -c "journalctl -f -u c2 | ncat --ssl logs4.papertrailapp.com 47837"
TimeoutStartSec=0
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
